# nives/rootfs/etc/s6-overlay/s6-rc.d/shodh/finish
#!/command/with-contenv bash
echo "[shodh] Process exited with code $1, restarting in 5s..."
sleep 5
