# nives/rootfs/etc/s6-overlay/s6-rc.d/init-options/run
#!/command/with-contenv bash
exec /usr/local/bin/options-to-env.sh
