# nives/rootfs/etc/s6-overlay/s6-rc.d/nives-server/finish
#!/command/with-contenv bash
echo "[server] Process exited with code $1, restarting in 3s..."
sleep 3
