# homemind-pro/rootfs/etc/s6-overlay/s6-rc.d/shodh/run
#!/command/with-contenv bash

SHODH_KEY=$(cat /data/.shodh_key 2>/dev/null || echo "")
if [ -z "$SHODH_KEY" ]; then
    echo "[shodh] ERROR: No API key found at /data/.shodh_key"
    sleep 5
    exit 1
fi

mkdir -p /data/shodh

echo "[shodh] Starting Shodh Memory on port 3030..."
export SHODH_API_KEYS="$SHODH_KEY"
export LD_LIBRARY_PATH="/usr/local/lib:${LD_LIBRARY_PATH:-}"
export SHODH_OFFLINE=true  # MiniLM-L6-v2 pre-downloaded into image; NER model still downloads on first boot (v0.2.0)
exec /usr/local/bin/shodh server \
    --host 127.0.0.1 \
    --port 3030 \
    --storage /data/shodh \
    --production
